# microservice1/apps/simple-pull/crontab
*/5 * * * * python /usr/src/app/pull-position.py >> results
